# TaskNode/src/escapewright.service
[Unit]
Description=EscapeWright Task Node
After=network.target

[Service]
Type=simple
Restart=on-failure
WorkingDirectory=/home/lodomo/EscapeWright/TaskNode/
ExecStart=/home/lodomo/EscapeWright/TaskNode/run.sh

[Install]
WantedBy=multi-user.target
