# TaskNode/src/escapewright.service
[Unit]
Description=EscapeWright Task Node
After=network.target

[Service]
Type=simple
StandardOutput=append:%hEscapeWright/TaskNode/logs/Output.log
StandardError=append:%hEscapeWright/TaskNode/logs/Error.log
Restart=on-failure
WorkingDirectory=%hEscapeWright/TaskNode/
ExecStart=./run.sh

[Install]
WantedBy=multi-user.target
