# TaskNode/src/escapewright.service
[Unit]
Description=EscapeWright Task Node
After=network.target

[Service]
Type=simple
StandardOutput=append:$HOME/EscapeWright/TaskNode/logs/Output.log
StandardError=append:$HOME/EscapeWright/TaskNode/logs/Error.log
Restart=on-failure
WorkingDirectory=$HOME/EscapeWright/TaskNode/
ExecStart=./run.sh

[Install]
WantedBy=multi-user.target
